# ping_st_changer.timer
[Unit]
Description=timer for ping_st_changer

[Timer]
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target